# systemd/fetch-chain.timer
[Unit]
Description=fetch CHAIN data

[Timer]
OnCalendar=*:21:00 UTC
Persistent=true

[Install]
WantedBy=kc2gprop.target
