# systemd/prop-storm.timer
[Unit]
Description=prop-storm every 3 hours

[Timer]
OnCalendar=0/3:07:00 UTC
Persistent=true

[Install]
WantedBy=kc2gprop.target
